# deploy/pi/holofan-update.service
# Pulls the latest code + image and (re)starts the stack. Run by holofan-update.timer.
# Because pull is a no-op when nothing changed, this doubles as a watchdog:
# if the app is down it comes back up; if a new version exists it updates.
[Unit]
Description=HoloFan Studio auto-update & watchdog
Requires=docker.service
After=docker.service network-online.target

[Service]
Type=oneshot
EnvironmentFile=-/etc/default/holofan
WorkingDirectory=/opt/holofan
# HOLOFAN_CLONE_DIR (set by install.sh) points at the git checkout, so updates git pull too.
ExecStart=/opt/holofan/holofan-update.sh ${HOLOFAN_CLONE_DIR}
